# files/ocean-data01/data01.mount.j2
[Unit]
Description=Mount 9p share data01 at /mnt/data01
After=network.target

[Mount]
What=data01
Where=/data01
Type=9p
Options=trans=virtio,version=9p2000.L

[Install]
WantedBy=multi-user.target