# files/cloudflare_ddns/cloudflare_ddns.timer.j2
[Unit]
Description={{ service }} timer

[Timer]
Unit={{ service }}.service
OnCalendar=*:00:*/5

[Install]
WantedBy=timers.target
