# files/cloudflare_ddns/cloudflare_ddns.timer.j2
[Unit]
Description={{ service }} timer

[Timer]
Unit={{ service }}.service
OnBootSec=1min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
