# ansible/files/crons/ntpdate.cron
0,11 * * * * root /usr/sbin/ntpdate time.nist.gov