# files/grafana-compose/grafana.service.j2
[Unit]
Description=Grafana Monitoring Dashboard with Docker Compose
Documentation=https://grafana.com/docs/
After=docker.service zfs-mount.service zfs-import.target network-online.target
Wants=network-online.target
Requires=docker.service zfs-mount.service zfs-import.target
RequiresMountsFor=/data01

[Service]
Type=oneshot
RemainAfterExit=yes
User=root
Group=root
WorkingDirectory={{ home }}

# Environment
Environment=COMPOSE_PROJECT_NAME={{ service }}
Environment=COMPOSE_HTTP_TIMEOUT=60

# Pre-start commands (idempotent)
ExecStartPre=/usr/bin/docker-compose -f {{ home }}/docker-compose.yml down --remove-orphans --timeout 30
ExecStartPre=-/usr/bin/docker-compose -f {{ home }}/docker-compose.yml pull --quiet

# Start command
ExecStart=/usr/bin/docker-compose -f {{ home }}/docker-compose.yml up -d

# Stop command  
ExecStop=/usr/bin/docker-compose -f {{ home }}/docker-compose.yml down

# Reload command - force recreation to pick up new config files
ExecReload=/usr/bin/docker-compose -f {{ home }}/docker-compose.yml down
ExecReload=/usr/bin/docker-compose -f {{ home }}/docker-compose.yml up -d --force-recreate

# Service behavior
Restart=on-failure
RestartSec=30
TimeoutStartSec=180
TimeoutStopSec=60

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ home }}
ReadOnlyPaths=/var/run/docker.sock

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier={{ service }}

[Install]
WantedBy=multi-user.target
