# homelab_playbooks/files/gitlab-docker/gitlab.service.j2
[Unit]
Description={{ service }}
After=docker.service
Requires=docker.service

[Service]
Type=simple
StandardOutput=journal
StandardError=journal
Restart=always
RestartSec=30s
TimeoutStartSec=10s
ExecStart=docker run --rm\
   --log-driver=journald \
   --name={{ service }} \
   --hostname={{ service }} \
   --shm-size 8gb \
   -e TZ="America/Los_Angeles" \
   -p {{ port_80 }}:80 \
   -p {{ port_443 }}:443 \
   -p {{ port_22 }}:22 \
   --volume {{ path_config }}/gitlab.rb:/etc/gitlab/gitlab.rb:ro \
   --volume {{ path_logs }}:/var/log/gitlab \
   --volume {{ path_data }}:/var/opt/gitlab \
  {{ image }}:{{ version }}
ExecStop=docker rm -f  {{ service }}

[Install]
WantedBy=multi-user.target
