# ansible/files/ocean-data01/data01.mount.j2
[Unit]
Description=ZFS Mount for /data01
After=zfs.target

[Mount]
What=/data01
Where=/data01
Type=zfs
Options=defaults

[Install]
WantedBy=multi-user.target