# files/homeassistant-compose/homeassistant.service.j2
[Unit]
Description=Home Assistant with Docker Compose
Documentation=https://www.home-assistant.io/docs/
# Start EARLY - immediately after docker to claim UDP ports before other host-networked services
After=docker.service zfs-mount.service zfs-import.target
Requires=docker.service zfs-mount.service zfs-import.target
RequiresMountsFor=/data01
# Other services should start AFTER this
Before=cloudflared.service node-exporter.service dcgm-exporter.service

[Service]
Type=oneshot
RemainAfterExit=yes
User={{ user }}
Group={{ user }}
WorkingDirectory={{ home }}

# Environment
Environment=COMPOSE_PROJECT_NAME={{ service }}
Environment=COMPOSE_HTTP_TIMEOUT=60

# Pre-start commands (idempotent)
ExecStartPre=/bin/bash -c 'cd {{ home }} && /usr/bin/docker compose down --remove-orphans --timeout 30'
ExecStartPre=-/bin/bash -c 'cd {{ home }} && /usr/bin/docker compose pull --quiet'

# Start command
ExecStart=/bin/bash -c 'cd {{ home }} && /usr/bin/docker compose up -d'

# Stop command
ExecStop=/bin/bash -c 'cd {{ home }} && /usr/bin/docker compose down'

# Reload command
ExecReload=/bin/bash -c 'cd {{ home }} && /usr/bin/docker compose restart'

# Service behavior
Restart=on-failure
RestartSec=30
TimeoutStartSec=300
TimeoutStopSec=60

# Security
NoNewPrivileges=true
ReadWritePaths={{ home }}
ReadOnlyPaths=/var/run/docker.sock

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier={{ service }}

[Install]
WantedBy=multi-user.target
