# files/agentbox/agentbox-remote-control.service.j2
[Unit]
Description=Claude Code remote-control (premium lane, Pro/Max subscription)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ user }}
Environment=HOME=/home/{{ user }}
Environment=PATH=/usr/local/bin:/usr/bin:/bin
# ANTHROPIC_API_KEY intentionally NOT set: Claude Code must use the subscription.
# An API key here would silently switch to metered billing.
WorkingDirectory={{ home }}/repos
ExecStart=/usr/bin/env claude remote-control --name agentbox
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agentbox-claude-rc

[Install]
WantedBy=multi-user.target
