# files/agentbox/agentbox-issue-watcher.timer.j2
[Unit]
Description=Run the agentbox issue watcher periodically

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
