# files/agentbox/agentbox-alert-receiver.service.j2
[Unit]
Description=agentbox alert receiver (Alertmanager -> ntfy + GitHub issues)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ user }}
Environment=HOME=/home/{{ user }}
Environment=PATH=/usr/local/bin:/usr/bin:/bin
# GH_TOKEN (for gh issue create) comes from the env file; no ANTHROPIC_API_KEY.
EnvironmentFile=-{{ home }}/.config/agentbox/agentbox.env
Environment=ALERT_NTFY_URL={{ alert_ntfy_url }}
Environment=ALERT_ISSUE_REPO={{ alert_issue_repo }}
Environment=ALERT_LISTEN_PORT={{ alert_receiver_port }}
ExecStart=/usr/bin/python3 /usr/local/bin/agentbox-alert-receiver.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agentbox-alert-receiver

[Install]
WantedBy=multi-user.target
