# files/agentbox/agentbox-rc@.service.j2
[Unit]
Description=Claude Code remote-control for %i (premium lane, Pro/Max subscription)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ user }}
Environment=HOME=/home/{{ user }}
Environment=PATH=/usr/local/bin:/usr/bin:/bin
# ANTHROPIC_API_KEY intentionally NOT set: Claude Code must use the subscription.
# An API key here would silently switch to metered billing.
# OTel env (telemetry endpoint + flags) without ANTHROPIC_API_KEY; per-repo label.
EnvironmentFile={{ home }}/.config/agentbox/agentbox.env
Environment=OTEL_RESOURCE_ATTRIBUTES=repo=%i,lane=claude,service=agentbox
WorkingDirectory={{ home }}/repos/%i
# Wrapper seeds workspace trust and answers the per-launch "Enable Remote
# Control?" prompt under a pty (Claude Code has no flag for either). See
# files/agentbox/rc-launch.sh.
ExecStart=/usr/local/bin/agentbox-rc-launch.sh %i
# always (not on-failure): the RC session can exit 0 on a claude.ai disconnect;
# bring it back regardless. An intentional `systemctl stop` is still respected.
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agentbox-rc-%i

[Install]
WantedBy=multi-user.target
