# roles/dns_infrastructure/templates/systemd/dns-stack.service.j2
[Unit]
Description=DNS Stack (PowerDNS + Kea DHCP + netboot.xyz)
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
TimeoutStartSec=120
TimeoutStopSec=60
WorkingDirectory={{ home }}
ExecStartPre=/usr/bin/docker compose pull --quiet
ExecStart=/usr/bin/docker compose up -d --remove-orphans
# Reclaim images left dangling by the :latest pull above. This host has no
# docker-cleanup timer (that lives in the standalone docker_ce.yaml, which the
# site pipeline does not run on dns hosts), so without this the old images pile
# up in /var/lib/docker and fill / over repeated deploys/reboots.
ExecStartPost=-/usr/bin/docker image prune -f
ExecStop=/usr/bin/docker compose down
ExecReload=/usr/bin/docker compose up -d --remove-orphans
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
